# scripts/linux/morfsync.service
# Service systemd pour morfSync.
#
# Le service tourne EN TANT QUE L'UTILISATEUR (User=/Group=), pour que ses
# données restent accessibles dans ~/.local/share/morfredus/morfSync
# (et non dans /var/lib, réservé au root). Le script install-service.sh
# remplace __RUN_USER__ et __RUN_HOME__ ; en installation manuelle, remplacez-les
# vous-même (utilisateur + son répertoire home).
#
# Nom d'unité : « morfsync ». Binaire et config dans __APP_DIR__
# (/opt/morfsync par défaut), comme tous les autres services du parc :
# un dossier FIXE hors du clone git, que déplacer le dépôt ne casse pas.
#
# Installation manuelle :
#   (préférer ./service.py install, qui fait tout et migre l'ancienne config)
#   sed -e 's/__RUN_USER__/votreuser/g' -e 's#__RUN_HOME__#/home/votreuser#g' \
#       scripts/linux/morfsync.service | sudo tee /etc/systemd/system/morfsync.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now morfsync

[Unit]
Description=morfSync — socle de synchronisation offline-first
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__RUN_USER__
Group=__RUN_USER__
# $HOME et le dossier de travail sont imposés explicitement : pour un service
# SYSTEM avec User=, systemd ne fournit pas toujours $HOME, et %h vaudrait /root.
# Sans ça, le hub ne résout pas son dossier de données et plante (SIGABRT).
Environment=HOME=__RUN_HOME__
WorkingDirectory=__RUN_HOME__
ExecStart=__APP_DIR__/morfSync __APP_DIR__/morfsync.json
Restart=on-failure
RestartSec=3
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
